# packaging/debian/prerm
#!/bin/sh

set -e

case "$1" in
    remove)
        # Stop service if running
        if systemctl is-active --quiet netctl; then
            systemctl stop netctl || true
        fi

        # Disable service
        if systemctl is-enabled --quiet netctl; then
            systemctl disable netctl || true
        fi

        # Clean up runtime files but preserve configuration
        if [ -d /var/run/netctl ]; then
            rm -rf /var/run/netctl/* 2>/dev/null || true
        fi

        echo "NetCtl service stopped. Configuration preserved in /etc/netctl"
        ;;

    upgrade|deconfigure)
        # Stop service during upgrade/deconfigure if active
        if systemctl is-active --quiet netctl; then
            systemctl stop netctl || true
        fi
        ;;

    *)
        echo "prerm called with unknown argument: $1" >&2
        exit 1
        ;;
esac

exit 0